# templates/tiaas.socket
[Unit]
Description=Galaxy tiaas service socket
[Socket]
ListenStream={{ tiaas_socket }}
[Install]
WantedBy=sockets.target
