# templates/tiaas.socket
[Unit]
Description=Galaxy tiaas service socket
[Socket]
ListenStream=/run/tiaas.sock
[Install]
WantedBy=sockets.target
